# com.gmail.invoices.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gmail.invoices.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string> <!-- Or wherever your python3 
is -->
    
<string>/Users/sivangalamidi/Development/Gmail/gmail_invoices.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>  <!-- Runs at 10PM -->
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/gmail_invoice_log.txt</string>
  <key>StandardErrorPath</key>
  <string>/tmp/gmail_invoice_error.txt</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>

